# supervisor.conf
[program:filebrowser]
command     = filebrowser -d /root/Ariang/filebrowser.db
directory   = /root/Ariang
user        = root
startsecs   = 3
